# crontab
0 10 * * * python /app/load_coin_csv.py 2>&1
