# crontab
* * * * * /usr/local/bin/python3 /app/load_coin_csv.py >> /tmp/load_coin_csv.log 2>&1
